# mcp-manager.service
[Unit]
Description=MCP Server Manager Daemon
# ネットワークが利用可能になってから起動する場合
# After=network.target

[Service]
# デーモンの実行ユーザー
# User=myuser

# Pythonデーモンスクリプトの絶対パス
# ExecStart=/usr/bin/python3 /path/to/your/project/mcp_manager.py
ExecStart=/usr/bin/python3 /data/systemd-mcp/mcp_manager.py

# デーモンスクリプトがあるディレクトリ
WorkingDirectory=/data/systemd-mcp

# systemdにプロセスを常に監視させ、異常終了した場合は再起動する
Restart=always
# 再起動するまでの待機時間
RestartSec=5

# 標準出力と標準エラー出力をsyslog (journald) に記録する
StandardOutput=journal
StandardError=journal

# 停止シグナルとしてSIGTERMを送信する
KillSignal=SIGTERM

[Install]
#どのターゲットで有効にするか
WantedBy=multi-user.target