# mcp-manager.service
[Unit]
Description=MCP Server Manager Daemon
# To start after network is available
# After=network.target

[Service]
# Daemon execution user
# User=myuser

# Absolute path to Python daemon script
# ExecStart=/usr/bin/python3 /path/to/your/project/mcp_manager.py
ExecStart=/usr/bin/python3 /usr/local/lib/mcp-manager/mcp_manager.py

# Let systemd always monitor the process and restart on abnormal termination
Restart=always
# Wait time before restart
RestartSec=5

# Record standard output and standard error to syslog (journald)
StandardOutput=journal
StandardError=journal

# Send SIGTERM as stop signal
KillSignal=SIGTERM

[Install]
# Which target to enable for
WantedBy=multi-user.target